# bkvy.service.example
[Unit]
Description=bkvy - LLM API Management Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=bkvy
Group=bkvy
WorkingDirectory=/opt/bkvy
Environment=PATH=/opt/bkvy/venv/bin
Environment=HOST=0.0.0.0
Environment=PORT=10006
Environment=LOG_LEVEL=info
Environment=WORKERS=1
# Uncomment to enable detailed transaction logging (CSV with full request details)
#Environment=TRANSACTION_LOGGING=true
# Uncomment to enable summary statistics (JSON with daily aggregates, lightweight)
#Environment=SUMMARY_STATS=true
# Dashboard configuration
# Uncomment to enable the dashboard (requires TRANSACTION_LOGGING=true)
#Environment=DASHBOARD_ENABLED=true
# Dashboard IP access control (comma-separated IPs, CIDR ranges, or wildcards)
# Examples: "127.0.0.1" (local only), "10.0.0.0/24" (network), "all" (no restrictions)
#Environment=DASHBOARD_ALLOWED_IPS=127.0.0.1
ExecStart=/opt/bkvy/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings (relaxed for compatibility)
NoNewPrivileges=yes

# Process limits
LimitNOFILE=1048576
LimitNPROC=1048576

[Install]
WantedBy=multi-user.target
